# system/taobao.cron
# 第1列分钟0~59
# 第2列小时0~23 (0表示子夜)
# 第3列日1~31
# 第4列月1~12
# 第5列星期0~7 (0和7都表示星期天)
# 第6列要运行的命令

# shop_list spider
0 23 * * * /usr/bin/task_shop_list

# goods_list spider
0 0 * * * /usr/bin/task_goods_list

# goods_detail spider
0 8 * * * /usr/bin/task_goods_detail
